# crontab
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app

# Test cron job - runs every hour
# * * * * * cd /app && /usr/local/bin/python scripts/test_cron.py >> /tmp/cron.log 2>&1
# 0 * * * * cd /app && /usr/local/bin/python scripts/test_cron.py >> /tmp/cron.log 2>&1

# Check youtube once an hour
0 * * * * cd /app && /usr/local/bin/python -m scripts.youtube.main >> /tmp/cron.log 2>&1

# Weekly GitHub sync on Mondays at 8 AM (commented out for testing)
# 0 8 * * 1 cd /app && /usr/local/bin/python scripts/github/repos_by_query.py >> /tmp/cron.log 2>&1